# deploy/hourly-editor.supervisord.conf
[program:hourly-editor]
directory=/var/www/predict/
user=www-data
autostart=true
autorestart=true
umask=002
command=venv/bin/uwsgi -y hourly/editor/uwsgi.yaml
