# deploy/hourly.supervisord.conf
[program:hourly]
directory=/var/www/predict/
user=www-data
autostart=true
autorestart=true
command=nice ionice -c3 venv/bin/python hourly/daemon.py
stopsignal=INT
stopwaitsecs=30
